# automation/macos/com.user.claude-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>SCRIPT_DIR/sync.sh</string>
        <string>push</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>HOME_DIR/.claude.json</string>
        <string>HOME_DIR/.claude/settings.json</string>
        <string>HOME_DIR/.claude/CLAUDE.md</string>
        <string>HOME_DIR/.claude/plugins/known_marketplaces.json</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>HOME_DIR/.claude-sync.log</string>

    <key>StandardErrorPath</key>
    <string>HOME_DIR/.claude-sync.log</string>
</dict>
</plist>
